# supervisord/conf.d/xfreerdp.conf
[program:freerdp]
command = xfreerdp /u:"%(ENV_RDP_USERNAME)s" /p:"%(ENV_RDP_PASSWORD)s" /v:"%(ENV_RDP_HOST)s":"%(ENV_RDP_PORT)s" %(ENV_RDP_ARGUMENTS)s
autorestart = true
